# .config/systemd/user/kanshi.service
[Unit]
Description=kanshi

[Service]
ExecStart=/home/shined/.local/bin/kanshi
Type=simple

[Install]
WantedBy=default.target
