# dot_config/systemd/user/kanshi.service
[Unit]
Description=kanshi

[Service]
ExecStart=kanshi
Type=simple

[Install]
WantedBy=default.target
